# recipes-rdm/wr-connect/wrc-startup-2.0/wrc_data.cron
#!/bin/sh

svc -o "/etc/daemontools/service/wrc_data"
exit 0
